# conf/mongodb_supervisor.ini
[program:startmongod]

command = sudo mongod -f /etc/mongod.conf
user = ubuntu
startsecs = 0

stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 5
stdout_logfile = /srv/quote_lookup/log/mongo.log
