# conf/quote_lookup_supervisor.ini
[program:quotelookup]

command     = /usr/local/bin/gunicorn -w4 -b 127.0.0.1:9000 manage:app
directory   = /srv/quote_lookup
user        = ubuntu
startsecs   = 3
startretries = 3  
stopwaitsecs = 0
autostart = true 
autorestart = true

redirect_stderr         = true
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = /srv/quote_lookup/log/app.log
stderr_logfile		= /srv/quote_lookup/log/err.log
